# ops/me.pdzeng.sync-writing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>me.pdzeng.sync-writing</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/panda/site/apps/pdzeng.com &amp;&amp; exec bun run publish:writing</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/panda/Library/Logs/pdzeng-writing-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/panda/Library/Logs/pdzeng-writing-sync.error.log</string>
</dict>
</plist>
